# scripts/blimas-monitor.service
[Unit]
Description=BLIMAS System Monitor
After=network.target mysql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/html/blimas
ExecStart=/usr/bin/php /var/www/html/blimas/scripts/system_monitor.php
Restart=always
RestartSec=300

[Install]
WantedBy=multi-user.target